# blind-controller.service
[Unit]
Description=Smart Blind Controller Service
After=network.target

[Service]
Type=simple
User=pi
# Adjust these paths to your installation
WorkingDirectory=/home/pi/blind-controller
ExecStart=/home/pi/blind-controller/venv/bin/python main.py
Restart=always
# Environment=MQTT_HOST=localhost

[Install]
WantedBy=multi-user.target
