# blind-controller.service
[Unit]
Description=Smart Blind Controller Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/blind-controller
ExecStart=/usr/bin/python3 main.py
Restart=always
Environment=MQTT_HOST=192.168.1.100

[Install]
WantedBy=multi-user.target
